# examples/systemd/em-tools.example.timer
# Copy to e.g. /etc/systemd/system/em-tools.timer
# Requires a service with the same basename: em-tools.service (see em-tools.example.service).
#
#   sudo cp examples/systemd/em-tools.example.service /etc/systemd/system/em-tools.service
#   sudo cp examples/systemd/em-tools.example.timer   /etc/systemd/system/em-tools.timer
# Edit paths/users and rake task in the .service file, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now em-tools.timer
#   systemctl list-timers em-tools.timer

[Unit]
Description=Periodic em-tools Rake run

[Timer]
# Daily at 03:30 (system timezone). See systemd.time(7) for syntax.
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
