# examples/systemd/em-tools.example.service
# Copy to e.g. /etc/systemd/system/em-tools.service
#   sudo cp examples/systemd/em-tools.example.service /etc/systemd/system/em-tools.service
# Edit User=, paths, and the rake task in ExecStart=, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now em-tools.service   # oneshot: runs once now
# Or enable the matching timer instead of the service directly.

[Unit]
Description=em-tools Rake (oneshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Replace with your login user (owns the repo and .env).
User=sky
Group=sky
WorkingDirectory=/home/sky/src/em-tools

# Optional: load KEY=value pairs from .env (same format as dotenv; no "export " prefix).
# Use absolute paths inside .env if you reference files.
EnvironmentFile=-/home/sky/src/em-tools/.env

# If bundle is not on systemd's PATH (rbenv/asdf), use login shell:
ExecStart=/usr/bin/bash -lc 'cd /home/sky/src/em-tools && bundle exec rake inventory:sync'

# Adjust timeout for long syncs.
TimeoutStartSec=7200
Nice=10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
