# schedule/systemd/em-tools-lowest-offer-snapshot.service.example
# em-tools — Amazon lowest-offer coverage snapshot (download seeds + publish).

[Unit]
Description=em-tools Amazon lowest-offer coverage snapshot
After=network-online.target em-tools-inventory-sync.service
Wants=network-online.target

[Service]
Type=oneshot
User=sky
Group=sky
WorkingDirectory=/home/sky/src/em-tools

EnvironmentFile=-/home/sky/src/em-tools/.env

ExecStart=/home/sky/src/em-tools/scripts/amazon-lowest-offer-snapshot.sh
Environment=EM_TOOLS_BUNDLE=/home/sky/.rbenv/shims/bundle

TimeoutStartSec=10800
Nice=10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
