# schedule/systemd/em-tools-inventory-sync.service.example
# em-tools — inventory sync (oneshot service, paired with .timer).
#
# Install:
#   sudo cp em-tools-inventory-sync.service.example /etc/systemd/system/em-tools-inventory-sync.service
#   sudo cp em-tools-inventory-sync.timer.example   /etc/systemd/system/em-tools-inventory-sync.timer
# Edit User=, paths, etc., then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now em-tools-inventory-sync.timer
#   journalctl -u em-tools-inventory-sync.service -e

[Unit]
Description=em-tools inventory sync (GCS CSV -> Elasticsearch)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=sky
Group=sky
WorkingDirectory=/home/sky/src/em-tools

# Source secrets and cluster URLs from the repo's .env.
EnvironmentFile=-/home/sky/src/em-tools/.env

# bash -lc picks up rbenv/asdf shims and other user-shell setup.
ExecStart=/usr/bin/bash -lc 'cd /home/sky/src/em-tools && bundle exec bin/em-tools inventory sync'

TimeoutStartSec=7200
Nice=10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
